# config/systemd/gist-sub.timer
[Unit]
Description=Run gist-sub fetch periodically

[Timer]
OnBootSec=30
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
